# south_carolina-scraper.service
[Unit]
Description=South Carolina Business Scraper - 24 Hour Scheduler
After=network.target

[Service]
Type=simple
User=your_username
WorkingDirectory=/path/to/bizleads-main
ExecStart=/usr/bin/python3 /path/to/bizleads-main/south_carolina_scheduler.py
Restart=always
RestartSec=60
StandardOutput=append:/var/log/south_carolina-scraper.log
StandardError=append:/var/log/south_carolina-scraper-error.log
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
